# examples/systemd/doc-audit-pull-mirrors.timer
[Unit]
Description=Refresh Doc Audit mirrors daily at 06:00 local time

[Timer]
OnCalendar=*-*-* 06:00:00
Persistent=true
Unit=doc-audit-pull-mirrors.service

[Install]
WantedBy=timers.target
